# IR/config/systemd/irremote.service
# /lib/systemd/system/irremote.service
# sudo systemctl enable irremote.service
[Unit]
Description=IR Remote Sender Service for Raspberry Pi
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/usr/local/bin/irremote_sender -m 192.168.1.100 -t home/ir/commands /home/pi/remotes/cambridge_cxa.json
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=HOME=/home/pi

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/pi

[Install]
WantedBy=multi-user.target
